# systemd/rag-query.service
[Unit]
Description=RAG Query API (search + stats) -> Chroma
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rag-split
Environment=LOG_LEVEL=INFO
Environment=CHROMA_HOST=127.0.0.1
Environment=CHROMA_PORT=7000
Environment=OLLAMA_EMBED_URL=http://127.0.0.1:11434/api/embeddings
Environment=EMBED_MODEL=nomic-embed-text
Environment=CHROMA_PATH=/mnt/ai-data/chroma

ExecStart=/opt/rag-split/venv/bin/uvicorn query_api:app --host 0.0.0.0 --port 8012
Restart=always
RestartSec=2
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
